# scripts/novnc.service
[Unit]
Description=noVNC Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Replaced by scripts/setup.sh at install time.
User=__NOVNC_USER__
# launch.sh is provided by the novnc package and wraps websockify startup reliably.
ExecStart=/usr/share/novnc/utils/launch.sh --listen 6080 --vnc localhost:5901
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
